# config/examples/rtl_433-supervisor.conf
[program:rtl_433]
command=/root/rtl_433/build/src/rtl_433 -R 123 -F 
"mqtt://loxberry-wohn:1883"user="loxberry",pass="MIFhegVK0HMPWIuj",events=BEER",user=loxberry
autostart=yes
autorestart=yes
startretries=100
stderr_logfile=REPLACELBPLOGDIR/rscan4lox.log
stdout_logfile=REPLACELBPLOGDIR/rscan4lox.log